# etc/cron.weekly/fstrim
#!/bin/sh

fstrim /boot/efi
fstrim /
fstrim /home
